# deploy/systemd/pushit-celery-beat.service
[Unit]
Description=PushIT Celery Beat Scheduler
After=network.target redis-server.service

[Service]
User=django
Group=www-data
WorkingDirectory=/var/www/django_websites/PushIT_server
ExecStart=/var/www/django_websites/PushIT_server/.venv/bin/celery -A config beat --loglevel=info --schedule=/var/www/django_websites/PushIT_server/celerybeat-schedule
Restart=on-failure
RestartSec=10
KillMode=mixed
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
