# deploy/systemd/pushit-web.service
[Unit]
Description=PushIT Gunicorn Web Server
After=network.target

[Service]
User=django
Group=www-data
WorkingDirectory=/var/www/django_websites/PushIT_server
RuntimeDirectory=pushit
ExecStart=/var/www/django_websites/PushIT_server/.venv/bin/gunicorn config.wsgi:application -c /var/www/django_websites/PushIT_server/deploy/gunicorn.conf.py
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5
KillMode=mixed
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
