# deploy/systemd/pushit-api-gunicorn.service
[Unit]
Description=PushIT API — Gunicorn (config.wsgi)
After=network.target pushit-env-fetch.service
# Hard dependency: without the SSM-fetched env, refuse to start.
Requires=pushit-env-fetch.service

[Service]
User=django
Group=www-data
WorkingDirectory=/var/www/django_websites/PushIT_server
EnvironmentFile=/run/pushit/.env
ExecStart=/var/www/django_websites/PushIT_server/.venv/bin/gunicorn config.wsgi:application -c /var/www/django_websites/PushIT_server/deploy/gunicorn.conf.py
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5
KillMode=mixed
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
